# sing-box/alpine/singbox
#!/sbin/openrc-run

description="sing-box Service"
command="/etc/sing-box/bin/sing-box"
command_args="run -c /etc/sing-box/config.json -C /etc/sing-box/conf"
pidfile="/run/sing-box.pid"

depend() {
    need net
    use logger
}

start_pre() {
    checkpath -f -o root:root /run/sing-box.pid
}

start() {
    ebegin "Starting sing-box"
    start-stop-daemon --start --background --make-pidfile --pidfile "$pidfile" --exec "$command" -- $command_args
    eend $?
}

stop() {
    ebegin "Stopping sing-box"
    start-stop-daemon --stop --pidfile "$pidfile" --retry 10
    eend $?
}